# supervisor/conf.d/tftpd.conf
[program:tftp]
command=in.tftpd -L --verbose -u root --secure --create /tftpboot
